# ops/poly-sniper.service
[Unit]
Description=Boreal Deterministic Trading Kernel (DTK v2)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/poly-trader/execution/dtk
# Run the highly optimized v2 kernel
ExecStart=/opt/poly-trader/execution/dtk/target/release/dtk
Restart=always
RestartSec=3

# Cap memory usage as a basic software HEAB safety mechanism
MemoryMax=1G

[Install]
WantedBy=multi-user.target
